# docker/open-tps-app/conf.d/opentps.conf
[program:opentps]
command=/app/opentps/start_opentps_linux.sh
autorestart=true