# docker/open-tps-app/conf.d/opentps.conf
[program:opentps]
command=opentps
autorestart=true